# Provisioning/RepositoryGateway/profile.d/global.settings
# Used to provide global GreenCom defaults; /etc/profile sources /etc/profile.d/*.sh
# See: https://help.ubuntu.com/community/EnvironmentVariables#System-wide_environment_variables
export GC_HOME=/home/pi/greencom
