# Provisioning/RepositoryGateway/cron.d/check_unsent_mongodb_values
45 05 * * * root /home/pi/greencom/bash/check_unsent_mongodb_values.sh >/dev/null 2>&1
