# Provisioning/RepositoryGateway/cron.d/check_disk_space
*/30 * * * * root /home/pi/greencom/bash/check_disk_space.sh >/dev/null 2>&1
