# Provisioning/RepositoryGateway/cron.d/ovpn_keepalive
*/20 * * * * root /home/pi/greencom/bash/keep_alive.sh >/dev/null 2>&1
