# Provisioning/RepositoryGateway/cron.d/py_post_alive
*/3 *	* * * 	root 	/home/pi/greencom/python/py-post-alive/post-alive.sh >/dev/null 2>&1
